# groups/default/etc/cron.d/simpleconf



*/2 * * * *    root    cd /var/cache/simpleconf/ ; git pull -q origin master &> /dev/null 
