# hosts/raspberrypi/etc/cron.d/rebooteveryfivedays
# Reboot every 5 days:                                                          
05 04 */5 * *   root    reboot 

